# scheduler.cron
* * * * * curl -s 'https://mde.currikistudio.org/api/api/v1/xapi-extract' >> /var/log/cron.log 2>&1
* * * * * php /var/www/html/artisan schedule:run >> /dev/null 2>&1

